# openclaw.service
[Unit]
Description=OCI A1 capacity watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/oracle-flex-create/acquire_a1.env
Environment=PATH=/root/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=/opt/oracle-flex-create
ExecStart=/usr/local/bin/acquire_a1.sh --env-file /etc/oracle-flex-create/acquire_a1.env
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
